# nots-dynamic-dns.service
[Unit]
Description=NOTS Dynamic DNS

[Service]
Type=simple
ExecStart=/opt/nots_dynamic_dns
User=YOURUSER
